# etc/rc.d/ffd
#!/bin/sh

daemon="/usr/local/bin/ffd"
daemon_flags="-d"
daemon_timeout=5

. /etc/rc.d/rc.subr

rc_reload=NO

rc_cmd $1
